# ops/oracle/window-db-backup.timer
[Unit]
Description=Nightly Window Neon backup export

[Timer]
OnCalendar=*-*-* 03:15:00 UTC
Persistent=true
Unit=window-db-backup.service

[Install]
WantedBy=timers.target
